# docker_cron/crontab
*/1 * * * * /bin/sh /root/scripts/hello_world.sh

